# ha_addon/rootfs/etc/services.d/scrypted-run/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Log Viewer
# Runs the Log Viewer
# ==============================================================================
declare name
declare value
bashio::log.info "Starting Scrypted..."

# Ensure the configuration exists
if bashio::fs.file_exists '/data/config/options.json'; then
    cp -f -v /data/config/options.json /etc/options
else
    mkdir -p -v /data/config \
        || bashio::exit.nok "Failed to create the Scrypted configuration directory"

    # Copy in template file
    cp -v /etc/options.json /data/config/
fi

for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting ${name} in scrypted run" 
    export "${name}=${value}"
done

if [[ ! -v SCRYPTED_ADMIN_USERNAME || -z "$SCRYPTED_ADMIN_USERNAME"  ]]; 
then
    bashio::log.info "No SCRYPTED_ADMIN_USERNAME setting found, skipping SCRYPTED_ADMIN_TOKEN env variable"
else
    bashio::log.info "SCRYPTED_ADMIN_USERNAME setting found, setting SCRYPTED_ADMIN_TOKEN env variable"
    export SCRYPTED_ADMIN_TOKEN=$(cat /tmp/.scrypted_admin_token)
fi

cd /server && npm --prefix /server exec scrypted-serve
