# ha_addon/rootfs/etc/services.d/nginx/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Community Add-on: Log Viewer
# Take down the S6 supervision tree when Nginx fails
# ==============================================================================
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services
