# ha_addon/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Scrypted
# Runs the Nginx daemon
# ==============================================================================

# Wait for the Scrypted service to become available
declare name
declare value
bashio::net.wait_for 10443

bashio::log.info "Starting NGinx..."
for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting ${name} in nginx run" 
    export "${name}=${value}"
done


if [[ ! -v SCRYPTED_ADMIN_USERNAME || -z "$SCRYPTED_ADMIN_USERNAME"  ]]; 
then
    bashio::log.info "No SCRYPTED_ADMIN_USERNAME setting found, removing proxy_pass_header config."
    sed -i "s/.*Authorization.*//g" /etc/nginx/servers/ingress.conf
else
    export SCRYPTED_ADMIN_TOKEN=$(echo $RANDOM$RANDOM$RANDOM$RANDOM$RANDOM | md5sum)
    bashio::log.info "SCRYPTED_ADMIN_USERNAME setting found, adding proxy_pass_header config." 
    sed -i "s/%%authorization_header%%/$SCRYPTED_ADMIN_TOKEN/g" /etc/nginx/servers/ingress.conf
fi


exec nginx
