# docker/fs/etc/s6-overlay/s6-rc.d/avahi/run
#!/command/execlineb -P

until [ -e /var/run/dbus/system_bus_socket ]; do
  echo "Waiting for dbus..."
  sleep 1s
done
echo "Starting Avahi daemon..."
exec avahi-daemon --no-chroot -f /etc/avahi/avahi-daemon.conf