# docker/fs/etc/s6-overlay/s6-rc.d/dbus/run
#!/bin/bash

echo "Starting dbus..."
exec dbus-daemon --system --nofork